# deployment/kleinanzeigen-scraper.service
[Unit]
Description=KleinanzeigenScraper Service
After=network.target

[Service]
Type=simple
# Use direct UID instead of username to avoid user lookup issues
User=0
Group=0
WorkingDirectory=/home/schneider/repos_private/KleinanzeigenScraper
ExecStart=/home/schneider/repos_private/KleinanzeigenScraper/scripts/start-service.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kleinanzeigen-scraper

[Install]
WantedBy=multi-user.target
